# config/deploy/shared/mqtt_subscriber_init.sh.erb
#!/bin/sh

### BEGIN INIT INFO
# Provides:          mqtt_subscriber
# Required-Start:    $local_fs $remote_fs $network $syslog
# Required-Stop:     $local_fs $remote_fs $network $syslog
# Default-Start:     2 3 4 5
# Default-Stop:      0 1 6
# Short-Description: starts the mqtt_subscriber 
# Description:       starts mqtt_subscriber
### END INIT INFO

set -e

# Feel free to change any of the following variables for your app:
TIMEOUT=${TIMEOUT-60}
APP_ROOT=<%= current_path %>
PID_DIR=$APP_ROOT/tmp/pids
PID=$PID_DIR/mqtt_subscriber.pid
CMD="cd $APP_ROOT;  bundle exec rake mqtt:sub RAILS_ENV=<%= fetch(:rails_env)  %> > $APP_ROOT/logs/mqtt_subscriber.log 2>&1"
AS_USER=<%= fetch(:deploy_user) %>
set -u

sig () {
  test -s "$PID" && kill -$1 `cat $PID`
}

run () {
  if [ "$(id -un)" = "$AS_USER" ]; then
    eval $1
  else
    su -c "$1" - $AS_USER
  fi
}

case "$1" in
start)
  sig 0 && echo >&2 "Already running" && exit 0
  run "$CMD"
  ;;
stop)
  sig QUIT && exit 0
  echo >&2 "Not running"
  ;;
force-stop)
  sig TERM && exit 0
  echo >&2 "Not running"
  ;;
restart|reload)
  sig QUIT
  run "$CMD"
  ;;
*)
  echo >&2 "Usage: $0 <start|stop|restart|force-stop>"
  exit 1
  ;;
esac
